# php/7.0-dev/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D

[program:fpm]
command=/usr/sbin/php-fpm7 -R --nodaemonize
